# sbase/.footprint
drwxr-xr-x	root/root	bin/
-rwxr-xr-x	root/root	bin/cat
-rwxr-xr-x	root/root	bin/chgrp
-rwxr-xr-x	root/root	bin/chmod
-rwxr-xr-x	root/root	bin/chown
-rwxr-xr-x	root/root	bin/cp
-rwxr-xr-x	root/root	bin/date
-rwxr-xr-x	root/root	bin/dd
-rwxr-xr-x	root/root	bin/echo
-rwxr-xr-x	root/root	bin/false
-rwxr-xr-x	root/root	bin/hostname
-rwxr-xr-x	root/root	bin/kill
-rwxr-xr-x	root/root	bin/ln
-rwxr-xr-x	root/root	bin/ls
-rwxr-xr-x	root/root	bin/mkdir
-rwxr-xr-x	root/root	bin/mknod
-rwxr-xr-x	root/root	bin/mv
-rwxr-xr-x	root/root	bin/pwd
-rwxr-xr-x	root/root	bin/readlink
-rwxr-xr-x	root/root	bin/rm
-rwxr-xr-x	root/root	bin/rmdir
-rwxr-xr-x	root/root	bin/sleep
-rwxr-xr-x	root/root	bin/sync
-rwxr-xr-x	root/root	bin/touch
-rwxr-xr-x	root/root	bin/true
-rwxr-xr-x	root/root	bin/uname
drwxr-xr-x	root/root	etc/
-rw-------	root/root	etc/crontab
drwxr-xr-x	root/root	etc/rc.d/
-rwxr-xr-x	root/root	etc/rc.d/cron
drwxr-xr-x	root/root	sbin/
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
lrwxrwxrwx	root/root	usr/bin/[ -> /usr/bin/test
-rwxr-xr-x	root/root	usr/bin/basename
-rwxr-xr-x	root/root	usr/bin/cal
-rwxr-xr-x	root/root	usr/bin/chroot
-rwxr-xr-x	root/root	usr/bin/cksum
-rwxr-xr-x	root/root	usr/bin/cols
-rwxr-xr-x	root/root	usr/bin/comm
-rwxr-xr-x	root/root	usr/bin/cron
-rwxr-xr-x	root/root	usr/bin/cut
-rwxr-xr-x	root/root	usr/bin/dirname
-rwxr-xr-x	root/root	usr/bin/du
-rwxr-xr-x	root/root	usr/bin/env
-rwxr-xr-x	root/root	usr/bin/expand
-rwxr-xr-x	root/root	usr/bin/expr
-rwxr-xr-x	root/root	usr/bin/flock
-rwxr-xr-x	root/root	usr/bin/fold
-rwxr-xr-x	root/root	usr/bin/head
-rwxr-xr-x	root/root	usr/bin/install
-rwxr-xr-x	root/root	usr/bin/join
-rwxr-xr-x	root/root	usr/bin/link
-rwxr-xr-x	root/root	usr/bin/logger
-rwxr-xr-x	root/root	usr/bin/logname
-rwxr-xr-x	root/root	usr/bin/md5sum
-rwxr-xr-x	root/root	usr/bin/mkfifo
-rwxr-xr-x	root/root	usr/bin/mktemp
-rwxr-xr-x	root/root	usr/bin/nice
-rwxr-xr-x	root/root	usr/bin/nl
-rwxr-xr-x	root/root	usr/bin/nohup
-rwxr-xr-x	root/root	usr/bin/od
-rwxr-xr-x	root/root	usr/bin/paste
-rwxr-xr-x	root/root	usr/bin/pathchk
-rwxr-xr-x	root/root	usr/bin/printenv
-rwxr-xr-x	root/root	usr/bin/printf
-rwxr-xr-x	root/root	usr/bin/renice
-rwxr-xr-x	root/root	usr/bin/rev
-rwxr-xr-x	root/root	usr/bin/runjobs
-rwxr-xr-x	root/root	usr/bin/seq
-rwxr-xr-x	root/root	usr/bin/setsid
-rwxr-xr-x	root/root	usr/bin/sha1sum
-rwxr-xr-x	root/root	usr/bin/sha224sum
-rwxr-xr-x	root/root	usr/bin/sha256sum
-rwxr-xr-x	root/root	usr/bin/sha384sum
-rwxr-xr-x	root/root	usr/bin/sha512-224sum
-rwxr-xr-x	root/root	usr/bin/sha512-256sum
-rwxr-xr-x	root/root	usr/bin/sha512sum
-rwxr-xr-x	root/root	usr/bin/sort
-rwxr-xr-x	root/root	usr/bin/split
-rwxr-xr-x	root/root	usr/bin/sponge
-rwxr-xr-x	root/root	usr/bin/tail
-rwxr-xr-x	root/root	usr/bin/tee
-rwxr-xr-x	root/root	usr/bin/test
-rwxr-xr-x	root/root	usr/bin/tftp
-rwxr-xr-x	root/root	usr/bin/time
-rwxr-xr-x	root/root	usr/bin/tr
-rwxr-xr-x	root/root	usr/bin/tsort
-rwxr-xr-x	root/root	usr/bin/tty
-rwxr-xr-x	root/root	usr/bin/unexpand
-rwxr-xr-x	root/root	usr/bin/uniq
-rwxr-xr-x	root/root	usr/bin/unlink
-rwxr-xr-x	root/root	usr/bin/uudecode
-rwxr-xr-x	root/root	usr/bin/uuencode
-rwxr-xr-x	root/root	usr/bin/wc
-rwxr-xr-x	root/root	usr/bin/which
-rwxr-xr-x	root/root	usr/bin/whoami
-rwxr-xr-x	root/root	usr/bin/yes
drwxr-xr-x	root/root	usr/sbin/
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/basename.1.gz
-rw-r--r--	root/root	usr/share/man/man1/cal.1.gz
-rw-r--r--	root/root	usr/share/man/man1/cat.1.gz
-rw-r--r--	root/root	usr/share/man/man1/chgrp.1.gz
-rw-r--r--	root/root	usr/share/man/man1/chmod.1.gz
-rw-r--r--	root/root	usr/share/man/man1/chown.1.gz
-rw-r--r--	root/root	usr/share/man/man1/chroot.1.gz
-rw-r--r--	root/root	usr/share/man/man1/cksum.1.gz
-rw-r--r--	root/root	usr/share/man/man1/cols.1.gz
-rw-r--r--	root/root	usr/share/man/man1/comm.1.gz
-rw-r--r--	root/root	usr/share/man/man1/cp.1.gz
-rw-r--r--	root/root	usr/share/man/man1/cron.1.gz
-rw-r--r--	root/root	usr/share/man/man1/cut.1.gz
-rw-r--r--	root/root	usr/share/man/man1/date.1.gz
-rw-r--r--	root/root	usr/share/man/man1/dd.1.gz
-rw-r--r--	root/root	usr/share/man/man1/dirname.1.gz
-rw-r--r--	root/root	usr/share/man/man1/du.1.gz
-rw-r--r--	root/root	usr/share/man/man1/echo.1.gz
-rw-r--r--	root/root	usr/share/man/man1/env.1.gz
-rw-r--r--	root/root	usr/share/man/man1/expand.1.gz
-rw-r--r--	root/root	usr/share/man/man1/expr.1.gz
-rw-r--r--	root/root	usr/share/man/man1/false.1.gz
-rw-r--r--	root/root	usr/share/man/man1/flock.1.gz
-rw-r--r--	root/root	usr/share/man/man1/fold.1.gz
-rw-r--r--	root/root	usr/share/man/man1/head.1.gz
-rw-r--r--	root/root	usr/share/man/man1/hostname.1.gz
-rw-r--r--	root/root	usr/share/man/man1/install.1.gz
-rw-r--r--	root/root	usr/share/man/man1/join.1.gz
-rw-r--r--	root/root	usr/share/man/man1/kill.1.gz
-rw-r--r--	root/root	usr/share/man/man1/link.1.gz
-rw-r--r--	root/root	usr/share/man/man1/ln.1.gz
-rw-r--r--	root/root	usr/share/man/man1/logger.1.gz
-rw-r--r--	root/root	usr/share/man/man1/logname.1.gz
-rw-r--r--	root/root	usr/share/man/man1/ls.1.gz
-rw-r--r--	root/root	usr/share/man/man1/md5sum.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mkdir.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mkfifo.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mknod.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mktemp.1.gz
-rw-r--r--	root/root	usr/share/man/man1/mv.1.gz
-rw-r--r--	root/root	usr/share/man/man1/nice.1.gz
-rw-r--r--	root/root	usr/share/man/man1/nl.1.gz
-rw-r--r--	root/root	usr/share/man/man1/nohup.1.gz
-rw-r--r--	root/root	usr/share/man/man1/od.1.gz
-rw-r--r--	root/root	usr/share/man/man1/paste.1.gz
-rw-r--r--	root/root	usr/share/man/man1/pathchk.1.gz
-rw-r--r--	root/root	usr/share/man/man1/printenv.1.gz
-rw-r--r--	root/root	usr/share/man/man1/printf.1.gz
-rw-r--r--	root/root	usr/share/man/man1/pwd.1.gz
-rw-r--r--	root/root	usr/share/man/man1/readlink.1.gz
-rw-r--r--	root/root	usr/share/man/man1/renice.1.gz
-rw-r--r--	root/root	usr/share/man/man1/rev.1.gz
-rw-r--r--	root/root	usr/share/man/man1/rm.1.gz
-rw-r--r--	root/root	usr/share/man/man1/rmdir.1.gz
-rw-r--r--	root/root	usr/share/man/man1/seq.1.gz
-rw-r--r--	root/root	usr/share/man/man1/setsid.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sha1sum.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sha224sum.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sha256sum.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sha384sum.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sha512-224sum.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sha512-256sum.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sha512sum.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sleep.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sort.1.gz
-rw-r--r--	root/root	usr/share/man/man1/split.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sponge.1.gz
-rw-r--r--	root/root	usr/share/man/man1/sync.1.gz
-rw-r--r--	root/root	usr/share/man/man1/tail.1.gz
-rw-r--r--	root/root	usr/share/man/man1/tee.1.gz
-rw-r--r--	root/root	usr/share/man/man1/test.1.gz
-rw-r--r--	root/root	usr/share/man/man1/tftp.1.gz
-rw-r--r--	root/root	usr/share/man/man1/time.1.gz
-rw-r--r--	root/root	usr/share/man/man1/touch.1.gz
-rw-r--r--	root/root	usr/share/man/man1/tr.1.gz
-rw-r--r--	root/root	usr/share/man/man1/true.1.gz
-rw-r--r--	root/root	usr/share/man/man1/tsort.1.gz
-rw-r--r--	root/root	usr/share/man/man1/tty.1.gz
-rw-r--r--	root/root	usr/share/man/man1/uname.1.gz
-rw-r--r--	root/root	usr/share/man/man1/unexpand.1.gz
-rw-r--r--	root/root	usr/share/man/man1/uniq.1.gz
-rw-r--r--	root/root	usr/share/man/man1/unlink.1.gz
-rw-r--r--	root/root	usr/share/man/man1/uudecode.1.gz
-rw-r--r--	root/root	usr/share/man/man1/uuencode.1.gz
-rw-r--r--	root/root	usr/share/man/man1/wc.1.gz
-rw-r--r--	root/root	usr/share/man/man1/which.1.gz
-rw-r--r--	root/root	usr/share/man/man1/whoami.1.gz
-rw-r--r--	root/root	usr/share/man/man1/yes.1.gz
drwxr-xr-x	root/root	usr/share/man/man8/
